# pressure/software/pressure.service
[Unit]
Description=Pressure daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/pressure_daemon.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
